# deploy/systemd/tdx-mqtt-listener.service
[Unit]
Description=TDX MQTT Listener - 即時通阻監聽 → Discord Webhook
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 /home/amakumo/.openclaw/workspace-nagi/skills/tdx-mqtt-listener/bin/tdx-mqtt-listener \
    --topics "v2/Bus/Alert/City/Kaohsiung,v2/Rail/Metro/Alert/KRTC,v3/Rail/TRA/Alert,v2/Rail/THSR/AlertInfo"
WorkingDirectory=/home/amakumo/.openclaw/workspace-nagi
EnvironmentFile=/home/amakumo/.openclaw/workspace/runtime/tdx/.env
Restart=always
RestartSec=30
RestartMaxDelaySec=300
RestartSteps=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
